# device/system-install/etc/systemd/system/can_setup.service
[Unit]
Description=CAN bus setup

[Service]
ExecStart=bash -x /usr/local/sbin/can_setup.sh

[Install]
WantedBy=multi-user.target
